# mail/.config/systemd/user/offlineimap.service
[Unit]
Description=Fetch mail from IMAP

[Service]
Type=oneshot
ExecStart=/bin/bash /home/crito/Applications/offlineimap.sh

[Install]
WantedBy=default.target
