# mail/systemd/index-mail.service
[Unit]
Description=Rebuild MU mail search index

[Service]
Type=oneshot
Nice=5
ExecStart=/usr/bin/mu index --maildir ~/.mail

[Install]
WantedBy=default.target
