# mail/systemd/offlineimap.service
[Unit]
Description=Fetch mail from IMAP

[Service]
Type=oneshot
Nice=5
ExecStart=/bin/bash %h/bin/offlineimap.sh

[Install]
WantedBy=default.target
