# mail/systemd/offlineimap.timer
[Unit]
Description=Check for mail every minute
RefuseManualStart=no
RefuseManualStop=no
 
[Timer]
Persistent=false
OnCalendar=minutely
 
[Install]
WantedBy=default.target
