# xorg/systemd/parcellite@.service
[Unit]
Description=handle clipboards

[Service]
Environment=DISPLAY=%I
ExecStart=/usr/bin/parcellite -d -n

[Install]
WantedBy=xorg@$i.target
