# xorg/systemd/xorg@.target
[Unit]
Description=Connection to Xorg (display %I)
